# coffeeapp.service
# Copyright (c) 2025 Cloudflight. All rights reserved. 
[Unit]
Description=coffee-app, presenting the slint.ui
After=weston.service
Requires=weston.service
 
[Service]
Restart=on-failure
Type=forking
Environment=WAYLAND_DISPLAY=/run/wayland-0 SLINT_FULLSCREEN=1
ExecStart=/usr/bin/NextCoffee
RestartSec=5

[Install]
WantedBy=multi-user.target
